# apps/api/ops/backup-daily.cron
# Executa todos os dias às 02:15 no servidor que aloja a API.
# O ambiente de deploy deve definir MONGODB_URI, STUDYFLOW_BACKUP_DIR e STUDYFLOW_BACKUP_RETENTION_DAYS.
15 2 * * * cd /srv/studyflow && npm --prefix real_dev/api run backup:daily >> /var/log/studyflow-backup.log 2>&1
